# .miniforge3/pkgs/conda-26.1.1-py313hfa70ccb_0/etc/profile.d/conda.csh
setenv CONDA_EXE "`cygpath 'D:\bld\conda_1772191268282\_h_env\Scripts\conda.exe'`";
setenv _CONDA_EXE "`cygpath 'D:\bld\conda_1772191268282\_h_env\Scripts\conda.exe'`";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "`cygpath 'D:\bld\conda_1772191268282\_h_env\python.exe'`";
setenv _CONDA_ROOT "`cygpath 'D:\bld\conda_1772191268282\_h_env'`";

source "`cygpath 'D:\bld\conda_1772191268282\_h_env\Lib\site-packages\conda\shell\etc\profile.d\conda.csh'`"