# .miniforge3/etc/profile.d/conda.csh
setenv CONDA_EXE "`cygpath 'C:/Users/manan/OneDrive2/Desktop/ashi/.miniforge3\Scripts\conda.exe'`";
setenv _CONDA_EXE "`cygpath 'C:/Users/manan/OneDrive2/Desktop/ashi/.miniforge3\Scripts\conda.exe'`";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "`cygpath 'C:/Users/manan/OneDrive2/Desktop/ashi/.miniforge3\python.exe'`";
setenv _CONDA_ROOT "`cygpath 'C:/Users/manan/OneDrive2/Desktop/ashi/.miniforge3'`";

source "`cygpath 'C:/Users/manan/OneDrive2/Desktop/ashi/.miniforge3\Lib\site-packages\conda\shell\etc\profile.d\conda.csh'`"